# deploy/docker/supervisord.conf
[supervisord]
nodaemon=true                 ; 在前台运行 supervisord
logfile=/dev/null             ; 将 supervisord 日志输出到标准输出/错误
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid  ; <-- [修正] 明确指定pidfile到可写的/tmp目录，解决权限问题

[program:redis]
command=/usr/bin/redis-server --loglevel notice
; user=appuser                  ; <-- [修正] 删除或注释掉此行。Docker容器已是非root用户，无需在此降权
autorestart=true
priority=10
stdout_logfile=/dev/stdout      ; 重定向 redis 标准输出到容器的标准输出
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr      ; 重定向 redis 标准错误到容器的标准错误
stderr_logfile_maxbytes=0

[program:gunicorn]
command=sh -c "echo '===== PRINTING ENV VARS =====' && printenv && echo '===== STARTING GUNICORN =====' && /usr/local/bin/gunicorn --bind 0.0.0.0:11235 --workers 1 --threads 4 --timeout 1800 --graceful-timeout 30 --keep-alive 300 --log-level info --worker-class uvicorn.workers.UvicornWorker server:app"
directory=/app                  ; 应用的工作目录
; user=appuser                  ; <-- [修正] 删除或注释掉此行。
autorestart=true
priority=20
environment=PYTHONUNBUFFERED=1  ; 确保 Python 输出直接发送到日志
stdout_logfile=/dev/stdout      ; 重定向 gunicorn 标准输出到容器的标准输出
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr      ; 重定向 gunicorn 标准错误到容器的标准错误
stderr_logfile_maxbytes=0

# Optional: Add filebeat or other logging agents here if needed
